# agent/com.focus.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.focus.agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/shrinilthakkar/git/focassist/.venv/bin/python</string>
        <string>-m</string>
        <string>agent.main</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/shrinilthakkar/git/focassist</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>FOCASSIST_BACKEND_URL</key>
        <string>http://13.205.122.54:8000</string>
        <key>FOCASSIST_TOKEN</key>
        <string>aedfb72a60b6c1c79bb31a3f3c31e5d577b6b4fa4a05ccb4babf38163ba36614</string>
        <key>AW_HOSTNAME</key>
        <string>SIMPL-ANALYTICS-016.local</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/focassist-agent.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/focassist-agent-error.log</string>

    <key>ThrottleInterval</key>
    <integer>300</integer>
</dict>
</plist>
